# worker-python/systemd/knowledgebaseservice-worker-python.service
# Knowledge Base Worker — systemd service
#
# Install:
#   sudo cp knowledgebaseservice-worker-python.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable knowledgebaseservice-worker-python.timer
#
# Manual run:
#   sudo systemctl start knowledgebaseservice-worker-python.service
#
# Logs:
#   journalctl -u knowledgebaseservice-worker-python.service

[Unit]
Description=Knowledge Base Service Worker
After=network.target postgresql.service

[Service]
Type=oneshot
User=limited_user
Group=limited_user
WorkingDirectory=/home/limited_user/applications/KnowledgeBaseService/worker-python
Environment="PATH=/home/limited_user/environments/knowledge-base-service/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/limited_user/applications/KnowledgeBaseService/worker-python/.env
ExecStart=/home/limited_user/environments/knowledge-base-service/bin/python main.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
