# worker-python/systemd/knowledgebaseservice-worker-python.timer
# Knowledge Base Worker — systemd timer
#
# Install:
#   sudo cp knowledgebaseservice-worker-python.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now knowledgebaseservice-worker-python.timer
#
# Check status:
#   systemctl list-timers --all | grep knowledgebaseservice

[Unit]
Description=Daily trigger for Knowledge Base Service Worker

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
